# sv/claymore/finish
#!/bin/bash
source ./miner.sh
exec runuser --login $MINER --command "tmux kill-session -t Claymore"
