# sv/claymore/run
#!/bin/bash
source ./miner.sh
exec runuser --login $MINER --command "/etc/sv/claymore/main.sh"
